# linux/systemd/system/mnt-truenas-Public.mount
[Unit]
Description=TrueNAS Public Share
After=network-online.target
Wants=network-online.target

[Mount]
What=//truenas/Public
Where=/mnt/truenas/Public
Type=cifs
Options=credentials=/home/william/.smbcredentials_truenas,uid=1000,gid=1000,file_mode=0755,dir_mode=0755,_netdev
TimeoutSec=30

[Install]
WantedBy=multi-user.target
