# linux/systemd/system/mnt-truenas-Public.automount
[Unit]
Description=Automount TrueNAS Public Share
ConditionPathExists=/mnt/truenas/Public

[Automount]
Where=/mnt/truenas/Public
TimeoutIdleSec=60

[Install]
WantedBy=multi-user.target
